# conf/supervisord.conf
[program:%(project)s]
command=%(django_user_home)s/scripts/rungunicorn_%(project)s.sh
directory=%(django_project_root)s
user=%(django_user)s
autostart=%(supervisor_autostart)s
autorestart=%(supervisor_autorestart)s
redirect_stderr=%(supervisor_redirect_stderr)s
stdout_logfile=%(supervisor_stdout_logfile)s